# docker/extraction/crontab

0 9 * * * / python3 usajobs/ingestion/app/engine/python/extractor/launchpad.py >> /var/log/cron.log 2>&1
